# apps/n8n/macos.n8n.trigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.odysea.n8n</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/curl</string>
      <string>-X</string>
      <string>POST</string>
      <!-- Replace with your actual n8n Production Webhook URL -->
      <string>http://localhost:8080/webhook/cd449724-2258-4f3a-833a-521338d32248</string>
    </array>

    <!-- Schedules the job to run daily at 00:00 (Midnight) -->
    <!-- If asleep at 00:00, it runs immediately on next wake -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>0</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/n8n_trigger.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/n8n_trigger_err.log</string>
  </dict>
</plist>